# web-build/cron.conf
[program:cron]
command=sudo /usr/sbin/cron -f
autorestart=true
startretries=10
