# web-build/time.cron.example
#ddev-generated
* * * * * date | tee -a /var/www/html/time.log
